# deploy/ubuntu/dante-main.service
[Unit]
Description=Dante quant factory — scanners & trading core (main.py)
Documentation=file:///home/ubuntu/dante_bots/Dual-Screener-Bot/RUNBOOK.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/dante_bots/Dual-Screener-Bot
EnvironmentFile=-/home/ubuntu/dante_bots/Dual-Screener-Bot/.env
ExecStart=/home/ubuntu/dante_bots/Dual-Screener-Bot/.venv/bin/python -u /home/ubuntu/dante_bots/Dual-Screener-Bot/main.py
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
