# deploy/systemd/dante-main.service.in
[Unit]
Description=Dante quant factory — scanners & trading core (shell entrypoint wrapper)
Documentation=file://@@INSTALL_ROOT@@/RUNBOOK.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=@@INSTALL_ROOT@@
EnvironmentFile=-@@INSTALL_ROOT@@/.env
ExecStart=/bin/bash @@INSTALL_ROOT@@/deploy/entrypoints/run_main_service.sh
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
