# deploy/systemd/dante-watchdog.service.in
[Unit]
Description=Dante DB heartbeat watchdog (oneshot, watchdog.py)
Documentation=file://@@INSTALL_ROOT@@/RUNBOOK.md
After=network-online.target dante-factory.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=DANTE_OPS_EVENTS_DB=@@INSTALL_ROOT@@/ops_events.sqlite
Environment=WATCHDOG_STATE_DIR=/var/lib/dante-watchdog
EnvironmentFile=-@@INSTALL_ROOT@@/.env
ExecStart=@@INSTALL_ROOT@@/.venv/bin/python -u @@INSTALL_ROOT@@/watchdog.py
StandardOutput=journal
StandardError=journal
