# deploy/systemd/dante-watchdog.timer
[Unit]
Description=Dante heartbeat watchdog timer (watchdog.py every 5 minutes)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=dante-watchdog.service

[Install]
WantedBy=timers.target
