# deploy/systemd/dante-factory.service.in
[Unit]
Description=Dante quant factory core (factory_launcher.py -> main.py)
Documentation=file://@@INSTALL_ROOT@@/RUNBOOK.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=@@INSTALL_ROOT@@
Environment=QUANT_SKIP_INLINE_TELEGRAM_DAEMON=1
Environment=DANTE_ASYNC_TELEGRAM_DAEMON=1
EnvironmentFile=-@@INSTALL_ROOT@@/.env
ExecStart=@@INSTALL_ROOT@@/.venv/bin/python -u @@INSTALL_ROOT@@/factory_launcher.py
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=120
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
