# deploy/systemd/dante-backup.timer
[Unit]
Description=Daily Dante DR backup timer (backup_to_cloud.sh)

[Timer]
# 서버 로컬 타임존 기준 매일 03:00 (RUNBOOK crontab 예와 동일). 변경 시 이 파일만 수정.
OnCalendar=*-*-* 03:00:00
Persistent=true
AccuracySec=30min
RandomizedDelaySec=1800
Unit=dante-backup.service

[Install]
WantedBy=timers.target
