# deploy/systemd/dante-dashboard.service.in
[Unit]
Description=Dante Streamlit dashboard
Documentation=file://@@INSTALL_ROOT@@/RUNBOOK.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=@@INSTALL_ROOT@@
EnvironmentFile=-@@INSTALL_ROOT@@/.env
ExecStart=@@INSTALL_ROOT@@/venv/bin/streamlit run @@INSTALL_ROOT@@/dashboard.py --server.address=0.0.0.0 --server.port=8501 --server.headless=true
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
