# deploy/systemd/dante-async.service.in
[Unit]
Description=Dante async Telegram queue consumer (async_telegram_daemon.py)
Documentation=file://@@INSTALL_ROOT@@/RUNBOOK.md
After=network-online.target dante-factory.service
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=@@INSTALL_ROOT@@
Environment=DANTE_ASYNC_TELEGRAM_DAEMON=1
EnvironmentFile=-@@INSTALL_ROOT@@/.env
ExecStart=@@INSTALL_ROOT@@/.venv/bin/python -u @@INSTALL_ROOT@@/async_telegram_daemon.py
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=90
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
