# deploy/systemd/dante-snapshot.service.in
[Unit]
Description=Dante CQRS DB snapshot (market_data_snapshot + ops_snapshot row)
Documentation=file://@@INSTALL_ROOT@@/RUNBOOK.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=2

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=@@INSTALL_ROOT@@
EnvironmentFile=-@@INSTALL_ROOT@@/.env
ExecStart=/bin/bash @@INSTALL_ROOT@@/deploy/entrypoints/run_snapshot_service.sh
StandardOutput=journal
StandardError=journal
