# deploy/factory.kr.crontab.example
# Dual-Screener-Bot — KR factory cron (→ /etc/cron.d/dual-screener-factory-kr)
#
# AUTO-GENERATED from factory_scan_schedule.py — do not edit by hand.
# Regenerate: python deploy/generate_factory_crontab.py
# Cycle-1 50 min from 10:00 KST; cycle-2 anchored before 15:30 close.
# install: sudo INSTALL_ROOT=... bash deploy/install_factory_cron.sh
#
# user/path: ubuntu · /home/ubuntu/dante_bots/Dual-Screener-Bot

SHELL=/bin/bash
CRON_TZ=Asia/Seoul
PATH=/usr/local/bin:/usr/bin:/bin

# --- KR staggered intraday (Mon–Fri KST, 8 slots; cycle-1 50 min) ---
0 10 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-supernova
50 10 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-nulrim
40 11 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-dante
30 12 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-ema5
20 13 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-master
10 14 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-bowl
0 15 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-supernova-r2
20 15 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --scan-kr-nulrim-r2

# --- KR post-close audit (well after intraday scans) ---
45 18 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --daily-kr

# --- US post-close audit (KST morning after NY close) ---
45 6 * * 2-6  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --daily-us

# --- Weekly Flow master (Saturday KST) ---
5 10 * * 6  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && TZ=Asia/Seoul bash ./factory.sh --weekly

# --- Legacy weekend brain surgery (optional) ---
0 10 * * 6  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && /home/ubuntu/dante_bots/Dual-Screener-Bot/venv/bin/python system_auto_pilot.py --run-autonomous-analysis-only >>/home/ubuntu/dante_bots/Dual-Screener-Bot/logs/autonomous_ubuntu.log 2>&1
