# deploy/factory.us.crontab.example
# Dual-Screener-Bot — US factory cron (→ /etc/cron.d/dual-screener-factory-us)
#
# AUTO-GENERATED from factory_scan_schedule.py — do not edit by hand.
# Regenerate: python deploy/generate_factory_crontab.py
# Staggered scans: 50 min apart, NY from 10:00 ET.
# CRON_TZ=America/New_York — DST 자동. Mon–Fri = US regular equity session days.
# install: sudo INSTALL_ROOT=... bash deploy/install_factory_cron.sh
#
# user/path: ubuntu · /home/ubuntu/dante_bots/Dual-Screener-Bot

SHELL=/bin/bash
CRON_TZ=America/New_York
PATH=/usr/local/bin:/usr/bin:/bin

# --- US staggered intraday (Mon–Fri ET, 9 slots, 50 min) ---
0 10 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-supernova
50 10 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-nulrim
40 11 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-dante
30 12 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-ema5
20 13 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-bowl
10 14 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-supernova-r2
0 15 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-nulrim-r2
50 15 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-dante-r2
40 16 * * 1-5  ubuntu  cd /home/ubuntu/dante_bots/Dual-Screener-Bot && bash ./factory.sh --scan-us-ema5-r2
