# bitget/deploy/systemd/dante-bitget-heatmap.service.in
[Unit]
Description=Bitget sector heatmap Streamlit (:8512)
Documentation=file://@@INSTALL_ROOT@@/bitget/RUNBOOK.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=@@INSTALL_ROOT@@
EnvironmentFile=-@@INSTALL_ROOT@@/.env
EnvironmentFile=-@@INSTALL_ROOT@@/bitget/.env
Environment=BITGET_HEATMAP_PORT=8512
ExecStart=/bin/bash @@INSTALL_ROOT@@/bitget/deploy/entrypoints/run_bitget_heatmap.sh
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=60
MemoryMax=512M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
